# clipboard-sync.service
[Unit]
Description=Clipboard Sync Daemon (Python version)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
ExecStart=/usr/bin/python3 %h/Documents/clipboard-sync-python/linux_sync.py
Restart=on-failure
RestartSec=5
Environment=DISPLAY=:0
Environment=WAYLAND_DISPLAY=wayland-0

[Install]
WantedBy=default.target
